# deploy/vpn-bot-admin-alerts.timer
[Unit]
Description=Run VPN bot admin alerts every 15 minutes

[Timer]
OnBootSec=5min
OnUnitActiveSec=15min
AccuracySec=1min

[Install]
WantedBy=timers.target
